# ansible/roles/aiserver/templates/nightly-tests.timer.j2
[Unit]
Description=Nightly Homelab Tests Timer

[Timer]
OnCalendar=*-*-* {{ nightly_tests_hour }}:{{ '%02d' | format(nightly_tests_minute) }}:00
Persistent=true

[Install]
WantedBy=timers.target
